# code/environments/production/modules/radicale/templates/radicale.service.epp
[Unit]
Description=Radicale CalDAV+CardDAV service
After=network.target

[Service]
Type=simple
ExecStart=<%= $radicale::install_dir %>/run.sh
WorkingDirectory=<%= $radicale::data_dir %>
User=<%= $radicale::user %>
Group=<%= $radicale::user %>
UMask=0027
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
